# playwright-mcp/supervisord.conf
; Supervisord config for playwright-mcp + VNC stack.
; Processes run as root (Playwright's chromium on Linux wants --no-sandbox
; anyway, which we pass via the mcp command; running as root inside a
; container is standard for this class of tool).

[supervisord]
nodaemon=false    ; top-level supervisord.conf already handles foreground
logfile=/tmp/supervisord.log

; 1. Virtual X display on :99, 1440x900 (comfortable for inspection).
[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1440x900x24 -nolisten tcp -ac
priority=100
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/tmp/xvfb.log
stderr_logfile=/tmp/xvfb.err

; 2. Minimal window manager so chromium windows have borders + can be focused.
;    (Chromium can run without one, but fluxbox makes the VNC view cleaner.)
[program:fluxbox]
command=/usr/bin/fluxbox
environment=DISPLAY=":99"
priority=200
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/tmp/fluxbox.log
stderr_logfile=/tmp/fluxbox.err

; 3. VNC server attached to the Xvfb display. -nopw because CF Access is
;    the real auth (this endpoint is only reachable via cloudflared +
;    email-gated CF Access app).
[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -forever -nopw -shared -rfbport 5900 -quiet
priority=300
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/tmp/x11vnc.log
stderr_logfile=/tmp/x11vnc.err

; 4. noVNC web UI (browser-based VNC client). Served on :8080.
;    Points to the x11vnc server on :5900. --web serves noVNC's static
;    assets (the HTML/JS client).
[program:novnc]
command=/usr/bin/websockify --web=/usr/share/novnc 8080 localhost:5900
priority=400
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/tmp/novnc.log
stderr_logfile=/tmp/novnc.err

; 5. playwright-mcp itself. Headed (no --headless), using DISPLAY=:99, so
;    every browser window it spawns renders on Xvfb and becomes visible
;    through x11vnc → noVNC.
;
;    --no-sandbox is required (chromium can't set up its sandbox inside
;    this container). --allowed-hosts guards against DNS rebinding.
[program:playwright-mcp]
command=/usr/bin/node /usr/lib/node_modules/@playwright/mcp/cli.js
  --browser chromium
  --no-sandbox
  --port 7004
  --host 0.0.0.0
  --allowed-hosts playwright-mcp,playwright-mcp:7004
  --output-dir /tmp/playwright-mcp
  --user-data-dir /data/browser
  --config /etc/playwright-mcp-config.json
environment=DISPLAY=":99"
priority=500
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/tmp/playwright-mcp.log
stderr_logfile=/tmp/playwright-mcp.err
